# docker/cron/crontab
* * * * * php /var/www/src/scripts/email.php >> /var/log/cron/cron.log 2>&1